# systemd/apt-news.service
[Unit]
Description=Update APT News

[Service]
Type=exec
ExecStart=/usr/bin/python3 /usr/lib/ubuntu-advantage/apt_news.py
